# platform/2-deployment/files/etc/systemd/system/evolver.service
[Unit]
Description=App server for evolver

[Service]
User=ufukty
WorkingDirectory=/home/ufukty
ExecStart=/home/ufukty/evolver
Restart=always

[Install]
WantedBy=multi-user.target