# platform/2-deployment/files/etc/systemd/system/runner.service
[Unit]
Description=App server for runner

[Service]
User=ufukty
WorkingDirectory=/home/ufukty
ExecStart=/home/ufukty/runner
Restart=always

[Install]
WantedBy=multi-user.target