# platform/stage/2-deployment/files/etc/systemd/system/runner.service
[Unit]
Description=App server for runner

[Service]
User=2iuFDs13YDedYc3N
WorkingDirectory=/home/2iuFDs13YDedYc3N
ExecStart=/home/2iuFDs13YDedYc3N/runner -config=/home/2iuFDs13YDedYc3N/config.yml
Restart=always

[Install]
WantedBy=multi-user.target