# platform/stage/image/internal/provisioner-files/map/etc/systemd/system/iptables-activation.service
[Unit]
Description=iptables rules for internal server image
Before=network-online.target
Wants=network-online.target

[Service]
Type=oneshot
ExecStart=/usr/sbin/iptables-restore /etc/iptables/iptables-rules.v4
ExecReload=/usr/sbin/iptables-restore /etc/iptables/iptables-rules.v4
RemainAfterExit=yes

[Install]
WantedBy=multi-user.target